# bdprochot.service
[Unit]
Description=Disable bdprochot

[Service]
ExecStart=/usr/bin/bdprochot.sh

[Install]
WantedBy=multi-user.target
